# deploy/systemd/maxclaw-gateway.service
[Unit]
Description=maxclaw Gateway
After=network-online.target maxclaw-bridge.service
Wants=network-online.target maxclaw-bridge.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__APP_DIR__
Environment=HOME=__HOME_DIR__
EnvironmentFile=-/etc/default/maxclaw
ExecStart=__APP_DIR__/scripts/run_gateway.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
