# deploy/systemd/nanobot-gateway.service
[Unit]
Description=nanobot-go Gateway
After=network-online.target nanobot-bridge.service
Wants=network-online.target nanobot-bridge.service

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__APP_DIR__
Environment=HOME=__HOME_DIR__
EnvironmentFile=-/etc/default/nanobot-go
ExecStart=__APP_DIR__/scripts/run_gateway.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
